# example/Krustacean-router.service
[Unit]
Description=Krustacean Router
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes

User=krab
Group=krab

ExecStartPre=/bin/bash -c '/usr/bin/grep -qE "^[0-9]+[[:space:]]+krustacean$" /etc/iproute2/rt_tables /etc/iproute2/rt_tables.d/*'
ExecStart=/usr/sbin/ip route replace local default dev lo table krustacean
ExecStart=/bin/bash -c '/usr/sbin/ip rule del fwmark 0x3001 table krustacean priority 100 2>/dev/null; /usr/sbin/ip rule add fwmark 0x3001 table krustacean priority 100'
ExecStop=/usr/sbin/ip rule del fwmark 0x3001 table krustacean priority 100
ExecStop=/usr/sbin/ip route del local default dev lo table krustacean

CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
NoNewPrivileges=yes

Restart=no

StandardOutput=journal
StandardError=journal

ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target